# deploy/verified-engineering-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Verified Engineering Studio
After=network-online.target ves.service
Wants=network-online.target ves.service

[Service]
Type=simple
User=pi
Group=pi
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --config /home/pi/verified-engineering-studio/deploy/cloudflared.yml run
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadOnlyPaths=/home/pi/verified-engineering-studio

[Install]
WantedBy=multi-user.target

