# deploy/ves.service
[Unit]
Description=Verified Engineering Studio public web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/verified-engineering-studio
EnvironmentFile=-/home/pi/verified-engineering-studio/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/verified-engineering-studio/.venv/bin/uvicorn ves.app:app --app-dir src --host 127.0.0.1 --port 8110
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/verified-engineering-studio

[Install]
WantedBy=multi-user.target

